# CONTAINER/citadel_napoleon.container
[Unit]
Description=CITADEL repo app napoleon (citadel)

[Container]
ContainerName=citadel-napoleon
HostName=citadel-napoleon
Image=localhost/citadel-napoleon:latest
Network=host
AddCapability=NET_ADMIN
AddCapability=NET_RAW
AddDevice=/dev/net/tun
EnvironmentFile=/home/openclaw/saf/NAPOLEON_HILLS_AI_MASTERMIND_CLASSES/CONTAINER/citadel_napoleon.env
Volume=citadel-napoleon-var-lib-tailscale:/var/lib/tailscale:Z
Volume=citadel-napoleon-opt-napoleon-hills-ai-mastermind-classes:/opt/NAPOLEON_HILLS_AI_MASTERMIND_CLASSES:Z
# PublishPort lines are ignored in host mode

[Service]
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
